# src/launchd/com.sutando.credential-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Sutando — credential-proxy launchd job (template).

  Role: keeps the quota-tracker credential-proxy alive at port 7846.
  The proxy caches ANTHROPIC_API_KEY reads so quota tracking doesn't
  call `claude config get` on every health-check interval (closes
  the port-7846-EADDRINUSE crash-loop described in issue #1086).

  Install:
    bash src/install-credential-proxy-launchd.sh

  Design:
    - KeepAlive=true so launchd restarts on crash (port becomes
      available once the dead process exits — no new EADDRINUSE).
    - ThrottleInterval=10: 10s minimum between restarts; prevents the
      3500-restart burst that burned CPU while the stale holder was
      alive.
    - Wrapper script kills any stale manual process on the port before
      exec'ing the real proxy, so the first boot after a manual run
      doesn't immediately fail with EADDRINUSE.
    - StandardOut/Err → workspace logs (not /tmp) so postmortems work.

  Placeholders substituted by install-credential-proxy-launchd.sh:
    __REPO__      — absolute path to the sutando repo checkout
    __WORKSPACE__ — absolute path to $SUTANDO_WORKSPACE
    __BREW_BIN__  — /opt/homebrew/bin or /usr/local/bin
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sutando.credential-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO__/src/launchd/credential-proxy-wrapper.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__BREW_BIN__:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>SUTANDO_WORKSPACE</key>
        <string>__WORKSPACE__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__WORKSPACE__/logs/credential-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>__WORKSPACE__/logs/credential-proxy.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
